# runit/uvicorn/run
#!/bin/sh

# 0E60835E-64B2-40C8-8865-3B5601EAC028

set -e -u -o pipefail

exec 2>&1

[ -f config ] && . $(realpath config)

: ${CHPST_ARGS=""}

if [ -d env ]; then
    CHPST_ARGS="${CHPST_ARGS} -e ${PWD}/env"
fi

ENVRC=""

if [ -f envrc ]; then
    ENVRC="$(cat envrc)"
fi

if [ -f envrc.local ]; then
    ENVRC="${ENVRC} $(cat envrc.local)"
fi

if [ -d envrc.d ] && [ $(ls envrc.d | wc -l) -gt 0 ]; then
    ENVRC="${ENVRC} $(cat envrc.d/*)"
fi

: ${UVICORN_USER:="nobody"}
: ${UVICORN_PATH:="/tmp"}
: ${UVICORN_ARGS=""}
: ${UVICORN_APP:="app:app"}
: ${UVICORN_VENV:="${UVICORN_PATH}/.venv"}
: ${UVICORN_PYTHON:="python3"}
: ${UVICORN_PYTHON_FLAGS=""}

cd ${UVICORN_PATH}

. ${UVICORN_VENV}/bin/activate

exec env -S "${ENVRC}" chpst -u ${UVICORN_USER} ${CHPST_ARGS} ${UVICORN_PYTHON} ${UVICORN_PYTHON_FLAGS} -m uvicorn ${UVICORN_ARGS} ${UVICORN_APP}
