# runit/gunicorn/run
#!/bin/sh

# E41C5C63-185F-43E8-A8E9-BF89A0BF0256

set -e -u

exec 2>&1

[ -f config ] && . $(realpath config)

: ${CHPST_FLAGS:=""}

if [ -d env ]; then
    CHPST_FLAGS="${CHPST_FLAGS} -e ${PWD}/env"
fi

ENVRC=""

if [ -f envrc ]; then
    ENVRC="$(cat envrc)"
fi

: ${GUNICORN_USER:="nobody"}
: ${GUNICORN_PATH:="/tmp"}
: ${GUNICORN_CMD_ARGS:=""}
: ${GUNICORN_APP:="app:app"}
: ${GUNICORN_VENV:="${GUNICORN_PATH}/.venv"}

cd ${GUNICORN_PATH}

. ${GUNICORN_VENV}/bin/activate

exec chpst -u ${GUNICORN_USER} ${CHPST_FLAGS} gunicorn ${GUNICORN_FLAGS} ${GUNICORN_APP}

# For uvicorn support: "--worker-class uvicorn.workers.UvicornWorker"
