# runit/gunicorn/run
#!/bin/sh

# e41c5c63-185f-43e8-a8e9-bf89a0bf0256

set -e -u -o pipefail

exec 2>&1

[ -f config ] && . $(realpath config)

: ${CHPST_ARGS=""}

if [ -d env ]; then
    CHPST_ARGS="${CHPST_ARGS} -e ${PWD}/env"
fi

ENVRC=""

if [ -f envrc ]; then
    ENVRC="$(cat envrc)"
fi

if [ -f envrc.local ]; then
    ENVRC="${ENVRC} $(cat envrc.local)"
fi

if [ -d envrc.d ] && [ $(ls envrc.d | wc -l) -gt 0 ]; then
    ENVRC="${ENVRC} $(cat envrc.d/*)"
fi

: ${GUNICORN_USER:="nobody"}
: ${GUNICORN_PATH:="/tmp"}
: ${GUNICORN_ARGS=""}
: ${GUNICORN_APP:="app:app"}
: ${GUNICORN_VENV:="${GUNICORN_PATH}/.venv"}
: ${GUNICORN_PYTHON:="python3"}
: ${GUNICORN_PYTHON_FLAGS=""}

cd ${GUNICORN_PATH}

. ${GUNICORN_VENV}/bin/activate

if [ "$(id -u)" = "0" ]; then
    CHPST_COMMAND="chpst -u ${GUNICORN_USER} ${CHPST_ARGS}"
fi

exec env -S "${ENVRC}" ${CHPST_COMMAND:-} ${GUNICORN_PYTHON} ${GUNICORN_PYTHON_FLAGS} -m gunicorn ${GUNICORN_ARGS} ${GUNICORN_APP}

# For logs: "--log-level debug"
# For uvicorn support: "--worker-class uvicorn.workers.UvicornWorker"
# For proxy address resolution: "--forwarded-allow-ips='*'"
