# runit/log/run
#!/bin/sh

set -e -u

mkdir -p main

exec chpst svlogd -ttt ./main
